# crontabs/fadein
# MIN HOUR DAY    MONTH DAYOFWEEK  COMMAND

  *   *    3,5    nov   thu        echo "3rd&Last Thursday"
  *   *    3,4    nov   thu        echo "3rd&4th Thursday"
  *   *    3      nov   thu        echo "3rd Thursday"
  *   *    1,5    *     tue-thu    echo "1st&last Tue-Thu"
 47  17     19    jan,nov thu      echo "5:45pm jan/nov 19th, on Thursday"


# vim: set tabstop=8:
