# crontabs/fadein
# MIN HOUR DAY MONTH DAYOFWEEK  COMMAND

    *     *   2     nov       sat  /usr/bin/echo "2nd Saturday of November"
    58    *   2     nov       sat  /usr/bin/echo "2nd Saturday of November 52 mins after"
    58    *   14     nov       *  /usr/bin/echo "14th of November 52 mins after"
    58    *   14     nov       *  /usr/bin/echo "14th of November 58 mins after"

# vim: set tabstop=8:
