# crontabs/fadein
# MIN HOUR DAY MONTH DAYOFWEEK  COMMAND

  *   *    3,5 nov   wed        /usr/bin/echo "3rd&Last Wednesday"
  *   *    3,4 nov   wed        /usr/bin/echo "3rd&4th Wednesday"
  *   *    3   nov   wed        /usr/bin/echo "3rd Wednesday"
  *   *    1,5 *     mon-wed    date

# vim: set tabstop=8:
